# launchagents/com.personalos.email-upload.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.personalos.email-upload</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/ryanhankins/personal-os/scripts/upload-email-report.sh</string>
  </array>
  <key>WatchPaths</key>
  <array>
    <string>/Users/ryanhankins/personal-os/data/last-email-report.json</string>
  </array>
  <key>StandardOutPath</key>
  <string>/Users/ryanhankins/personal-os/logs/email-upload-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ryanhankins/personal-os/logs/email-upload-stderr.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
